# com.user.pdf-to-md.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service label (must be unique) -->
    <key>Label</key>
    <string>com.user.pdf-to-md</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <!-- REPLACE_PYTHON_PATH: This will be replaced by setup.sh -->
        <string>REPLACE_PYTHON_PATH</string>
        <!-- REPLACE_SCRIPT_PATH: This will be replaced by setup.sh -->
        <string>REPLACE_SCRIPT_PATH</string>
    </array>

    <!-- Run at load (system startup/user login) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep the daemon running (restart if it crashes) -->
    <key>KeepAlive</key>
    <true/>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <!-- REPLACE_WORKING_DIR: This will be replaced by setup.sh -->
    <string>REPLACE_WORKING_DIR</string>

    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/tmp/pdf-to-md.out.log</string>

    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>/tmp/pdf-to-md.err.log</string>

    <!-- Run as user (not root) -->
    <key>UserName</key>
    <string>REPLACE_USERNAME</string>

    <!-- Environment variables (if needed) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>

